# dist/dev.closedmesh.closedmesh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Reference launchd LaunchAgent for Senda.

  The official installer (`curl senda.network/install | sh -s -- --service`)
  generates this file with the right paths for your home directory at
  ~/Library/LaunchAgents/network.senda.runtime.plist. This copy ships in
  the release tarball as a reference.

  Replace USERNAME below with your macOS username if you copy this manually.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>network.senda.runtime</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/USERNAME/.local/bin/senda</string>
        <string>serve</string>
        <string>--auto</string>
        <string>--mesh-name</string>
        <string>senda</string>
        <string>--join-url</string>
        <string>https://entry.senda.network/api/status</string>
        <string>--headless</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/USERNAME</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ProcessType</key>
    <string>Background</string>
    <key>StandardOutPath</key>
    <string>/Users/USERNAME/Library/Logs/senda/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/USERNAME/Library/Logs/senda/stderr.log</string>
</dict>
</plist>
