# dist/network.senda.runtime.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd LaunchAgent template for Senda.

  The official installer (`curl senda.network/install | sh -s -- --service`)
  renders this template into ~/Library/LaunchAgents/network.senda.runtime.plist
  with the right paths for your home directory. The @UPPERCASE@ tokens are
  substituted at install time; this copy ships in the release tarball as a
  reference for what the installed agent looks like.

  Manage it after install with: senda service start|stop|status|logs
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>@SERVICE_LABEL@</string>
    <key>ProgramArguments</key>
    <array>
        <string>@SERVICE_RUNNER@</string>
    </array>
    <key>WorkingDirectory</key>
    <string>@HOME_DIR@</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ProcessType</key>
    <string>Background</string>
    <key>StandardOutPath</key>
    <string>@STDOUT_LOG@</string>
    <key>StandardErrorPath</key>
    <string>@STDERR_LOG@</string>
</dict>
</plist>
